# scripts/root.cron
*/5 * * * *					/data/nginxconf/scripts/nginx-HUP-from-container.sh
