# scripts/root.cron
*/5 * * * *					/data/nginx/scripts/nginx-HUP-from-container.sh
